# snapper-snapshots.path
# File: snapper-snapshots.path
# Project Path: ./snapper-snapshots.path
# Installation Path: /usr/lib/systemd/system/snapper-snapshots.path
#
# Systemd path file to monitor Snapper snapshots

[Unit]
Description=Monitors for new Snapper snapshots

[Path]
PathModified=/.snapshots
Unit=update-systemd-boot-snapshots.service

[Install]
WantedBy=multi-user.target
