# timeshift-snapshots.path
# File: timeshift-snapshots.path
# Project Path: ./timeshift-snapshots.path
# Installation Path: /usr/lib/systemd/system/timeshift-snapshots.path
#
# Systemd path file to monitor Timeshift snapshots

[Unit]
Description=Monitors for new Timeshift snapshots

[Path]
PathModified=/run/timeshift/backup/timeshift-btrfs/snapshots
Unit=update-systemd-boot-snapshots.service

[Install]
WantedBy=multi-user.target
