# src/systemd-boot-snapshots-0.1.0/systemd-boot-entries.path
# File: systemd-boot-entries.path
# Project Path: ./systemd-boot-entries.path
# Installation Path: /usr/lib/systemd/system/systemd-boot-entries.path
#
# Systemd path file to monitor bootloader entries

[Unit]
Description=Monitors for new boot loader entries
DefaultDependencies=no
Requires=boot-efi.mount
After=boot-efi.mount
BindsTo=boot-efi.mount

[Path]
PathModified=/boot/efi/loader/entries
Unit=update-systemd-boot-snapshots.service

[Install]
WantedBy=boot-efi.mount
